# packages/django-eremaea2/django-eremaea2-2.0.16.tar.gz/django-eremaea2-2.0.16/eremaea-pull@.service
[Unit]
Description=Pull django-eremaea2 for %i
After=network.target
PartOf=eremaea.target

[Service]
Type=simple
ExecStart=/usr/bin/eremaeactl $EREMAEACTL_GLOBAL_OPTS pull $EREMAEACTL_PULL_OPTS
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal
LimitNPROC=1
LimitFSIZE=0

[Install]
WantedBy=multi-user.target eremaea.target
