# packages/django-eremaea2/django-eremaea2-2.0.16.tar.gz/django-eremaea2-2.0.16/eremaea.target
[Unit]
Description=A target allowing to start/stop all eremaea-pull@.service and eremaea-purge@.timer instances at once

