# packages/django-eremaea2/django-eremaea2-2.0.16.tar.gz/django-eremaea2-2.0.16/eremaea-purge@.timer
[Unit]
Description=Purge django-eremaea2 for %i once a hour
After=network.target
PartOf=eremaea.target

[Timer]
OnCalendar=hourly
AccuracySec=5min

[Install]
WantedBy=multi-user.target eremaea.target
